# deploy/scan-worker.service
# Skills Threat Analysis — systemd Service
#
# 安装:
#   sudo cp deploy/scan-worker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable scan-worker
#   sudo systemctl start scan-worker
#
# 管理:
#   sudo systemctl status scan-worker
#   sudo systemctl restart scan-worker
#   sudo journalctl -u scan-worker -f

[Unit]
Description=Skills Threat Analysis Worker
After=network.target

[Service]
Type=simple
User=chailinyuan
Group=chailinyuan
WorkingDirectory=/home/chailinyuan/skills-threat-analysis
Environment=PYTHONPATH=/home/chailinyuan/skills-threat-analysis/src
ExecStart=/home/chailinyuan/miniconda3/envs/skill/bin/python -m scanner.worker.cli --config config.yaml --workers 4
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=60

StandardOutput=append:/home/chailinyuan/skills-threat-analysis/logs/worker.log
StandardError=append:/home/chailinyuan/skills-threat-analysis/logs/worker.log

[Install]
WantedBy=multi-user.target
