# deploy/supervisord.conf
; Skills Threat Analysis — Supervisord Configuration
;
; 安装:  pip install supervisor
; 启动:  supervisord -c deploy/supervisord.conf
; 管理:  supervisorctl -c deploy/supervisord.conf status|start|stop|restart scan-worker
; 日志:  tail -f /home/chailinyuan/skills-threat-analysis/logs/worker-stdout.log

[unix_http_server]
file=/tmp/scan-worker-supervisor.sock

[supervisord]
logfile=/home/chailinyuan/skills-threat-analysis/logs/supervisord.log
pidfile=/tmp/scan-worker-supervisord.pid
directory=/home/chailinyuan/skills-threat-analysis

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/scan-worker-supervisor.sock

[program:scan-worker]
command=/home/chailinyuan/miniconda3/envs/skill/bin/python -m scanner.worker.cli --config config.yaml --workers 4
directory=/home/chailinyuan/skills-threat-analysis
environment=PYTHONPATH="/home/chailinyuan/skills-threat-analysis/src"
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=60
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/home/chailinyuan/skills-threat-analysis/logs/worker-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
